# docker/root/etc/s6/web/run
#!/usr/bin/env bash

set -e

cd /app

gosu shotsmart python api.py

